# phase2/process-manager/systemd/cineworld-otp.service
[Unit]
Description=CineWorld OTP Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cineworld/services/otp
EnvironmentFile=/opt/cineworld/config/.env
ExecStart=/opt/cineworld/services/otp/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8002
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cineworld-otp

[Install]
WantedBy=multi-user.target
