# scripts/systemd/cineworld-booking.service
[Unit]
Description=CineWorld Booking Service
After=network.target mysql.service cineworld-identity.service cineworld-catalog.service
Wants=mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cineworld/services/booking
EnvironmentFile=/opt/cineworld/config/.env
ExecStart=/opt/cineworld/services/booking/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8004
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cineworld-booking

[Install]
WantedBy=multi-user.target
