# scripts/systemd/cineworld-scheduler.service
[Unit]
Description=CineWorld Scheduler Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cineworld/services/scheduler
EnvironmentFile=/opt/cineworld/config/.env
ExecStart=/opt/cineworld/services/scheduler/venv/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cineworld-scheduler

[Install]
WantedBy=multi-user.target
